# scripts/kubelet_restarter.timer
[Unit]
Description=Kubelet Restarter

[Timer]
OnUnitActiveSec=60s
OnBootSec=60s

[Install]
WantedBy=timers.target