# scripts/kubelet_restarter.service
[Unit]
Description=Kubelet Restarter
After=kubelet.service
[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/bash /root/kubelet_restarter.sh